# deploy/tarco-backend.service
[Unit]
Description=Tarco FastAPI backend
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=kevin
Group=kevin
WorkingDirectory=/home/kevin/tarco/bin/backend
Environment=TARCO_HOME=/home/kevin/tarco
Environment=PATH=/home/kevin/tarco/bin/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=LOG_DIR=/home/kevin/tarco/log
Environment=HOST=127.0.0.1
Environment=PORT=8000
EnvironmentFile=-/home/kevin/tarco/conf/backend.env
ExecStart=/home/kevin/tarco/bin/tarco-backend
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
