# deploy/tarco-backend.service
[Unit]
Description=Tarco FastAPI backend
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=kevin
Group=kevin
WorkingDirectory=/mnt/mac/Users/liujianan/Develop/tarco/backend
Environment=PATH=/home/kevin/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOST=127.0.0.1
Environment=PORT=8000
ExecStart=/home/kevin/.local/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
